# examples/systemd.mount
[Unit]
Description=Buildkit Ramdisk

[Mount]
What=tmpfs
Where=%%RAMDISK%%
Type=tmpfs
Options=nosuid,size=8G

[Install]
WantedBy=multi-user.target
